# edgebal/supervisord.conf
[supervisord]
nodaemon=true

[program:balagent]
command=/home/ambergreen/balagent

[program:nginx]
command=nginx -g 'daemon off;'
